# etc/cron.d/minereu
#update system time. every hour
0 *     * * *     root   /usr/sbin/ntpdate -u 0.uk.pool.ntp.org

#watch dog constantly check the cgminer instance to make sure it is always running.
#every 30 seconds
* *     * * *     root  /opt/scripta/startup/miner-start.sh 2>&1 |/usr/bin/logger -t minereu_watchdog
* *     * * *     root  sleep 30;/opt/scripta/startup/miner-start.sh 2>&1 |/usr/bin/logger -t minereu_watchdog

#report and hashrate update
*/5 * * * * www-data /usr/bin/php5 /opt/scripta/etc/cron.d/5min/ALERTdeviceCount 2>&1 |/usr/bin/logger -t scripta_background
*/5 * * * * www-data /usr/bin/php5 /opt/scripta/etc/cron.d/5min/ALERThashrate 2>&1 |/usr/bin/logger -t scripta_background
*/5 * * * * www-data /usr/bin/php5 /opt/scripta/etc/cron.d/5min/hashrate 2>&1 |/usr/bin/logger -t scripta_background
*/5 * * * * www-data /usr/bin/php5 /opt/scripta/etc/cron.d/5min/temp 2>&1 |/usr/bin/logger -t scripta_background
#
